# crontab
30 2 * * * /usr/local/bin/python /app/aurora.py >> /var/log/cron.log 2>&1
